# ops/web/supervisord.conf
[supervisord]
nodaemon = true

[program:nginx]
command = nginx -c /app/ops/web/nginx.conf
startsecs = 5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:app-gunicorn]
command = gunicorn -w4 --log-level=debug --chdir=/app -b 127.0.0.1:5000 app:app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
